# packaging/systemd/linkup-agent.service
[Unit]
Description=Linkup agent — API/état (Laravel)
Documentation=https://example.invalid/linkup
After=network.target linkup-bridge.service
Wants=linkup-bridge.service

[Service]
Type=simple
WorkingDirectory=/opt/linkup/agent
# Écoute 0.0.0.0 : le téléphone attaque l'agent sur le LAN (handshake pairing).
# post/upload_max_filesize relevés : l'envoi PC→tél (dashboard) poste le fichier
# en multipart à Laravel ; 8 Mo (défaut PHP) rejetait les vidéos.
ExecStart=/usr/bin/php -d post_max_size=256M -d upload_max_filesize=256M artisan serve --host=0.0.0.0 --port=8000
Restart=on-failure
RestartSec=2
EnvironmentFile=-%h/.config/linkup/agent.env

[Install]
WantedBy=default.target
