# infra/videohub/linkup-videohub.conf
; Supervisor — service VideoHub Linkup (téléchargeur vidéo + transcript).
; À déposer dans /etc/supervisor/conf.d/ sur le VPS.
[program:linkup-videohub]
directory=/var/www/projects/linkup/videohub
command=/var/www/projects/linkup/videohub/.venv/bin/python run.py
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/linkup/videohub.err.log
stdout_logfile=/var/log/linkup/videohub.out.log
; PATH inclut /usr/local/bin pour que yt-dlp trouve `deno` (solveur EJS du défi JS
; YouTube). DENO_DIR pointe vers un cache inscriptible par www-data (pas de HOME).
; Le POT provider bgutil tourne en conteneur Docker sur 127.0.0.1:4416.
environment=PYTHONUNBUFFERED="1",PATH="/usr/local/bin:/usr/bin:/bin",DENO_DIR="/var/www/projects/linkup/videohub/.deno-cache"
user=www-data
